# rootfs/etc/s6-overlay/s6-rc.d/skystats/run
#!/command/with-contenv sh
cd /app/core
exec s6wrap --prepend=skystats --timestamps --args /app/core/skystats
